# packaging/deb/systemd/prometheus-labels-db-recorder.timer
[Unit]
Description=Run prometheus-labels-db-recorder periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
